# github/com.paper-agent.github-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paper-agent.github-sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/lunadong/Code/paper_agent/paper_collection && /usr/bin/python3 daily_update.py --check --no-email && cd /Users/lunadong/Code/paper_agent/github && /usr/bin/python3 github_push.py</string>
    </array>
    
    <!-- Run daily at 5pm -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Also run when the computer wakes from sleep (catch-up) -->
    <key>StartOnMount</key>
    <false/>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/paper-github-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/paper-github-sync.error.log</string>
    
    <!-- Keep alive settings - run again if machine was asleep during scheduled time -->
    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/Users/lunadong/Code/paper_agent</key>
            <true/>
        </dict>
    </dict>
</dict>
</plist>
